# laml/scripts/launchd/com.laml.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.laml.dashboard</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>REPO_ROOT/laml/scripts/start-laml-dashboard.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>REPO_ROOT/laml/dashboard</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>REPO_HOME/Library/Logs/LAML/dashboard.log</string>
	<key>StandardErrorPath</key>
	<string>REPO_HOME/Library/Logs/LAML/dashboard.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
